# deployment/systemd/csky-api.service
[Unit]
Description=CartoSky API
After=network.target

[Service]
Type=simple
User=cartosky
Group=cartosky
UMask=0027
WorkingDirectory=/opt/cartosky
EnvironmentFile=-/etc/cartosky/api.env
Environment=PATH=/opt/cartosky/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/opt/cartosky/.venv/bin/uvicorn backend.app.main:app --host 127.0.0.1 --port 8200
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target